# assets/myserver_start.conf
[program:myserver]
priority = 168
command = /usr/local/bin/server -u -c /etc/myserver/myserver.conf 
startsecs = 0
autorestart = true
